# units/netns@.service.in
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/usr/sbin ip netns | grep %I || /usr/sbin/ip netns add %I"
ExecStop=/usr/sbin/ip netns del %I
